# control-plane/macos/LaunchAgents/com.jarvis.ai.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Name -->
    <key>Label</key>
    <string>com.jarvis.ai</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/YOUR_USERNAME/Jarvis/dist/main.js</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/Jarvis</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
    </dict>

    <!-- Run at load (start on login) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive - restart if crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Restart after N seconds if exits -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Jarvis/logs/jarvis-stdout.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Jarvis/logs/jarvis-stderr.log</string>

    <!-- Nice level (process priority) -->
    <key>Nice</key>
    <integer>0</integer>

    <!-- Process type -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Disabled by default (enable with: launchctl load) -->
    <key>Disabled</key>
    <false/>

    <!-- Session type -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Exit timeout (seconds to wait before SIGKILL) -->
    <key>ExitTimeOut</key>
    <integer>30</integer>

    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>10240</integer>
    </dict>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>10240</integer>
    </dict>
</dict>
</plist>
